# infra/ollama/systemd/ollama.service.d/idr-dify.conf
# Dify(Podman)가 host.containers.internal(호스트 LAN IP)로 접속할 수 있도록
# Ollama를 루프백이 아닌 모든 인터페이스에 바인딩합니다.
#
# 적용: 프로젝트 루트에서 ./infra/ollama/apply-dify-host-bind.sh
[Service]
Environment="OLLAMA_HOST=0.0.0.0:11434"
